# infraestrutura/trabalho_agendado.cron
*/5 * * * * cd /home/ubuntu/infraestrutura && /usr/bin/sudo /usr/bin/docker compose -f /home/ubuntu/infraestrutura/docker-compose.yml run --rm jar >> /var/log/processar_planilha.log 2>&1

